# homeserver-unlock.service
[Unit]
Description=%p service

[Service]
#Environment=DEBUG=x
ExecStart=/home/florian/git/cryptroot-unlocker/cryptroot-unlocker.sh
Type=oneshot

